# roles/shipping/files/shipping.service
[Unit]
Description=Shipping Service
After=network.target

[Service]
User=roboshop
Environment="CART_ENDPOINT=cart-dev.yourtherapist.in:8080"
Environment="DB_HOST=mongodb-dev.yourtherapist.in"
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping
Restart=always

[Install]
WantedBy=multi-user.target
